# examples/contextlake.timer
# Example systemd user timer — runs the refresh every 30 minutes (and shortly
# after login). Pair with contextlake.service. See that file for install steps.

[Unit]
Description=Run contextlake refresh periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
